# mnt.service
[Unit]
Description=mnt
After=network.target
[Service]
ExecStart=home/rd/mnt.sh
[Install]
WantedBy=default.target
